# api/estetik_api.service
[Unit]
Description=Runs the Estetik API.
[Service]
ExecStart=/srv/estetik/api/run_app.sh
[Install]
RequiredBy=multi-user.target